# pinger.service
[Unit]
Description=Pinger keepalive for Render URLs
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
# Путь к клону репозитория (при установке в домашнюю папку смените на /home/USER/pinger и User=USER)
WorkingDirectory=/opt/pinger
Environment="KEEPALIVE_LOGFILE=/var/log/keepalive_random.log"
ExecStart=/usr/bin/python3 -u keepalive_random.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pinger

[Install]
WantedBy=multi-user.target
